# backend/deploy/supervisord.conf
[supervisord]
nodaemon=true

[program:api]
directory=/opt/binance-auto-trader
command=/opt/binance-auto-trader/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
autostart=true
autorestart=true
stdout_logfile=/var/log/binance-auto-trader/api.out.log
stderr_logfile=/var/log/binance-auto-trader/api.err.log
